# dyndns.service
[Unit]
Description=Update External IP address to a Dynamic DNS server
Wants=network-online.target
After=network-online.target

[Service]
User=#USER#
# Set the following to the source program directory (where you git cloned)
WorkingDirectory=#PROGDIR#
ExecStart=#PROGDIR#/%p

[Install]
WantedBy=default.target
